# tvoc-logger.service.example
[Unit]
Description=TVOC Logger Service
# ネットワーク接続が確立されるまで待機します
After=network-online.target
Wants=network-online.target

[Service]
# スクリプトを実行するユーザーとグループを指定します
User=pi
Group=pi

# スクリプトが存在するディレクトリのフルパスを指定します
# .envファイルなどを正しく読み込むために重要です
WorkingDirectory=/home/pi/tvoc

# サービスの実行前に少し待機します。ハードウェアの初期化待ちに有効です。
ExecStartPre=/bin/sleep 15

# 実行するコマンドをフルパスで指定します
# 仮想環境(.venv)内のPythonを使用するようにしてください
ExecStart=/home/pi/tvoc/.venv/bin/python3 /home/pi/tvoc/log_tvoc.py

# エラーで終了した場合、常に10秒後に再起動します
Restart=always
RestartSec=10

[Install]
# システム起動時にサービスが有効になるようにします
WantedBy=multi-user.target